# service_config/etc-uwsgi.ini
[program:etc-uwsgi]
command=uwsgi --ini /app/uwsgi.ini
directory=/app
autostart = true     ; 在 supervisord 启动的时候也自动启动
stopsignal=QUIT
user=root
startsecs = 5     ; 启动 5 秒后没有异常退出，就当作已经正常启动了
startretries = 3   ; 启动失败自动重试次数，默认是 3
autorestart = true   ; 程序异常退出后自动重启
redirect_stderr = true  ; 把 stderr 重定向到 stdout，默认 false
stdout_logfile_maxbytes = 10KB  ; stdout 日志文件大小，默认 50MB
disable-logging = False
stdout_logfile = /app/logs/my_service/supervisor_log/api_stdout.log
stderr_logfile = /app/logs/my_service/supervisor_log/api_err.log
